# startup_scripts/TeleInfod
#!/sbin/runscript
# Copyright Laurent Faillie 2015
#
# Start/Stop TeleInfo daemon on OpenRC system (tested on Gentoo)

depend() {
	need net
	after mosquitto uart
}

TID_BIN=${DMKDEM_BIN:-/usr/local/sbin/TeleInfod}
TID_PID=${DMKDEM_PID:-/var/run/TeleInfod.pid}
TID_USR=${DMKDEM_USR:-teleinfod:nobody}
# use following command to create this dedicated user
# useradd --comment "TeleInfod injector user" --gid nogroup --groups uucp --system --no-user-group --shell /sbin/nologin --home-dir /dev/null teleinfod

start() {
	ebegin "Starting TeleInfod daemon"
	start-stop-daemon --start --exec ${TID_BIN} \
		--interpreted --user ${TID_USR} \
		--background --make-pidfile --pidfile ${TID_PID}
	eend $?
}

stop() {
	ebegin "Stopping TeleInfod daemon"
	start-stop-daemon --stop --exec ${TID_BIN} \
		--pidfile ${TID_PID}
	eend $?
}

